# deploy/launchd/com.autoresearch.daily_update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily data update + classify — 09:15 IST (NSE opens 09:15; ingest runs
  in parallel with the pre-market scan and finishes before the 10:00 IST
  execution window).

  Extends prices/macro, pulls 24h news, refreshes the earnings calendar.
  Idempotent (upsert / LLM-cache short-circuit) — safe to re-run.

  LLM classification is SKIPPED here (--skip-classify; also the code
  default as of 2026-05-18): the committed momentum-quality strategy
  consumes only prices/universe/sectors, so daily classify is wasted
  cost and would add noise to the 4-week paper-validation logs. To
  re-enable, drop the flag / pass --no-skip-classify.

  Install:
    1. Confirm WorkingDirectory + the uv path below match this machine.
    2. cp this file to ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.autoresearch.daily_update.plist
    4. Verify: launchctl list | grep daily_update

  Logs to state/logs/daily_update.{out,err}.log. Machine clock is IST.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.autoresearch.daily_update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>cd "/Users/aryanmehta/Desktop/My Work/autoresearch-trading-india" &amp;&amp; /opt/homebrew/bin/uv run python -m scripts.daily_update --skip-classify</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/aryanmehta/Desktop/My Work/autoresearch-trading-india/state/logs/daily_update.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/aryanmehta/Desktop/My Work/autoresearch-trading-india/state/logs/daily_update.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
