# deploy/launchd/com.autoresearch.daily_report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Post-close risk check + daily report — 15:35 IST (NSE closes 15:30).

  Summarises the day (positions, ledger, risk signals) and refreshes the
  dashboard. Read-only over the ledger; never places orders.

  Install:
    1. Confirm WorkingDirectory + the uv path below match this machine.
    2. cp this file to ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.autoresearch.daily_report.plist
    4. Verify: launchctl list | grep daily_report

  Logs to state/logs/daily_report.{out,err}.log. Machine clock is IST.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.autoresearch.daily_report</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>cd "/Users/aryanmehta/Desktop/My Work/autoresearch-trading-india" &amp;&amp; /opt/homebrew/bin/uv run python -m scripts.daily_report</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>15</integer>
    <key>Minute</key>
    <integer>35</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/aryanmehta/Desktop/My Work/autoresearch-trading-india/state/logs/daily_report.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/aryanmehta/Desktop/My Work/autoresearch-trading-india/state/logs/daily_report.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
